# GV25MonitoringApp/PythonWorker/supervisord.conf
[supervisord]
nodaemon=true

[program:aa]
command=/ # instert first job

[program:bb]
command=/ # insert second job