# GV25MonitoringApp/PythonWorker/supervisord.conf
[supervisord]
nodaemon=true

[program:FakeDataGen]
command=python ./scripts/fakeDataGen.py # instert first job

[program:bb]
command=/ # insert second job